# templates/etc/init.d/carbon-aggregator.erb
#!/bin/sh
# !!! generated by puppet !!!

### BEGIN INIT INFO
# chkconfig: 2345 99 80
# Provides: carbon-aggregator
# Required-Start: $remote_fs $network $time
# Required-Stop: $remote_fs $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description: Start/Stop carbon-aggregator
# Description: Enables Graphites carbon-aggregator data aggregation engine
### END INIT INFO

PYTHON_CMD='python -W ignore'
CARBON_CONF='/opt/graphite/conf/carbon.conf'

OPERATION="$1"
if [ $# -gt 1 ]; then
    shift
    INSTANCES=$@
else
    INSTANCES=$(awk '$1 ~ /^\[aggregator/ { sub("\[aggregator:?","", $1); sub("]$","", $1); if ($1 == "") { print "a"} else { print $1} }' "${CARBON_CONF}")
fi

case "${OPERATION}" in
    start)
        for instance in ${INSTANCES}; do
            ${PYTHON_CMD} /opt/graphite/bin/carbon-aggregator.py --instance=${instance} start
        done
        ;;
    stop)
        for instance in ${INSTANCES}; do
            ${PYTHON_CMD} /opt/graphite/bin/carbon-aggregator.py --instance=${instance} stop
        done
        ;;
    status)
        for instance in ${INSTANCES}; do
            ${PYTHON_CMD} /opt/graphite/bin/carbon-aggregator.py --instance=${instance} status
        done
        ;;
    restart)
        $0 stop ; sleep 3 ; $0 start
        RC=$?
        exit $RC
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status} [instance instance ...]"
        exit 1
        ;;
esac

exit
