# templates/etc/systemd/gunicorn.service.erb
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
PIDFile=/run/gunicorn.pid
WorkingDirectory=/opt/graphite/webapp/graphite
ExecStart=/usr/bin/gunicorn --pid /run/gunicorn.pid --timeout=<%= scope.lookupvar('graphite::gunicorn_arg_timeout') %> --bind=<%= scope.lookupvar('graphite::gunicorn_bind') %> --workers=<%= scope.lookupvar('graphite::gunicorn_workers') %> --user <%= scope.lookupvar('graphite::gr_web_user_REAL') %> --group <%= scope.lookupvar('graphite::gr_web_group_REAL') %> --access-logfile /opt/graphite/storage/access-gunicorn.log --error-logfile /opt/graphite/storage/error-gunicorn.log wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
